# com.lancehunt.dotfiles-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lancehunt.dotfiles-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/lhunt/code/lancehunt/dotfiles/bin/sync.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/lhunt/code/lancehunt/dotfiles/sync.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/lhunt/code/lancehunt/dotfiles/sync.error.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
